# base/etc/systemd/system/relmon.service
[Unit]
Description=relmon

[Service]
Type=oneshot
User=rybalkin
Group=rybalkin
ExecStart=/usr/bin/python /usr/bin/relmon --update
