# base/etc/systemd/system/srcfetcher.service
[Unit]
Description=srcfetcher

[Service]
Type=oneshot
User=rybalkin
Group=rybalkin
ExecStart=/usr/bin/python /usr/bin/srcfetcher --quiet
