# base/etc/systemd/system/syncema.service
[Unit]
Description=syncema

[Service]
Type=oneshot
User=rybalkin
Group=rybalkin
ExecStart=/usr/bin/python -m mr.syncema
